# systemd/dot-config/systemd/user/borg-backup@.service
[Unit]
Description=Borg backup service for %i

[Service]
ExecStart=%h/bin/autogenerated-%i.backup.sh
Type=oneshot
