# systemd/dot-config/systemd/user/borg-backup@.service
[Unit]
Description=Borg backup service for %i

[Service]
Slice=backup.slice
MemoryAccounting=yes
IOAccounting=yes
Nice=15
IOSchedulingClass=idle
IOWeight=50
ExecStart=/usr/bin/flock %t/borg-backup.lock %h/bin/autogenerated-%i.backup.sh
Type=oneshot
