# systemd/dot-config/systemd/user/dev1-config-files.service
[Unit]
Description=Borg dev1-config-files backup service
Wants=dev1-config-files.timer

[Service]
ExecStart=%h/bin/autogenerated-dev1-config-files.backup.sh
Type=oneshot
RemainAfterExit=yes
