# systemd/dot-config/systemd/user/dev2-PKM.service
[Unit]
Description=Borg dev2-PKM backup service
Wants=dev2-PKM.timer

[Service]
ExecStart=%h/bin/autogenerated-dev2-PKM.backup.sh
Type=oneshot
RemainAfterExit=yes
