# systemd/dot-config/systemd/user/dev2-config-files.service
[Unit]
Description=Borg dev2-config-files backup service
Wants=dev2-config-files.timer

[Service]
ExecStart=%h/bin/autogenerated-dev2-config-files.backup.sh
Type=oneshot
