# systemd/dot-config/systemd/user/dev3-config-files.service
[Unit]
Description=Borg dev3-config-files backup service
Wants=dev3-config-files.timer

[Service]
ExecStart=%h/bin/autogenerated-dev3-config-files.backup.sh
Type=oneshot
RemainAfterExit=yes
