# systemd/dot-config/systemd/user/dev3-PKM.service
[Unit]
Description=Borg dev3-PKM backup service
Wants=dev3-PKM.timer

[Service]
ExecStart=%h/bin/autogenerated-dev3-PKM.backup.sh
Type=oneshot
RemainAfterExit=yes
