# systemd/dot-config/systemd/user/feh.service
[Unit]
Description=Apply desktop wallpaper with feh
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/bin/apply-wallpaper

[Install]
WantedBy=graphical-session.target
