# PKM/dot-config/systemd/user/browserexport.service
[Unit]
Description=Browser export
Wants=browserexport.timer

[Service]
Type=oneshot
ExecStart=%h/bin/backup-browser.sh
RemainAfterExit=yes
