# PKM/dot-config/systemd/user/foam-crypt-auto.service
[Unit]
Description=Guarded Foam rclone crypt bisync
After=syncthing.service

[Service]
Type=oneshot
Environment=FOAM_CRYPT_DEBOUNCE=8
ExecStart=%h/bin/foam-crypt-auto run
