# PKM/dot-config/systemd/user/hypexport.service
[Unit]
Description=Hypexport
Wants=hypexport.timer

[Service]
Type=oneshot
ExecStart=%h/bin/backup-hypothesis.sh
