# PKM/dot-config/systemd/user/PKM-export@.service
[Unit]
Description=PKM %i export service

[Service]
Type=oneshot
ExecStart=%h/bin/%i.sh
