# PKM/dot-config/systemd/user/promnesia-index.service
[Unit]
Description=Promnesia index
Wants=promnesia-index.service

[Service]
Type=oneshot
ExecStart=%h/bin/promnesia index
RemainAfterExit=yes
