# utilities/dot-config/systemd/user/display-health-check.service
[Unit]
Description=Check display scanout health
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/bin/display-health-check
Environment="DISPLAY=:0"
Environment="XAUTHORITY=%h/.Xauthority"
