# utilities/dot-config/systemd/user/dotfiles-spikes.service
[Unit]
Description=Update Foam system spike report
After=graphical-session.target

[Service]
Type=oneshot
Nice=15
IOSchedulingClass=idle
ExecStart=%h/bin/dotfiles-spikes update
