# utilities/dot-config/systemd/user/setup-displays.service
[Unit]
Description=Setup displays

[Service]
ExecStart=%h/bin/setup-displays.sh
Environment="DISPLAY=:0"
Environment="XAUTHORITY=%h/.Xauthority"
