# assets/other/backup.slice
[Unit]
Description=Background backup workload slice

[Slice]
MemoryAccounting=yes
IOAccounting=yes
CPUQuota=40%
MemoryHigh=1G
MemoryMax=2G
