# assets/other/borg-root-backup@.service
[Unit]
Description=Borg root backup service for %i

[Service]
ExecStart=/usr/bin/autogenerated-%i.backup.sh
User=root
Group=root
Type=oneshot
