# assets/other/borg-root-backup@.service
[Unit]
Description=Borg root backup service for %i

[Service]
Slice=backup.slice
MemoryAccounting=yes
IOAccounting=yes
Nice=15
IOSchedulingClass=idle
IOWeight=50
ExecStart=/usr/bin/flock %t/borg-root-backup.lock /usr/bin/autogenerated-%i.backup.sh
User=root
Group=root
Type=oneshot
