# peripherals/dot-config/systemd/user/cursor_indicator@.service
[Unit]
Description=cursor indicator

[Service]
ExecStart=%h/bin/cursor_indicator %i
