# peripherals/dot-config/systemd/user/keyd-observer.service
[Unit]
Description=Keyd observer
After=plasma-workspace.target

[Service]
ExecStart=%h/bin/keyd-observer

[Install]
WantedBy=plasma-workspace.target
