# audio/dot-config/systemd/user/configure-scarlett-monitoring.service
[Unit]
Description=Configure Scarlett 4i4 hardware monitoring
After=pipewire.service wireplumber.service
Wants=pipewire.service wireplumber.service

[Service]
Type=oneshot
ExecStart=%h/.local/bin/configure-scarlett-monitoring

[Install]
WantedBy=default.target
